# examples/haproxy/services/syslogd/run
#!/bin/sh

exec syslogd -n -O /dev/stdout
